# scripts/launchd/com.carlosblanco.monitor-subastas-valencia.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.carlosblanco.monitor-subastas-valencia.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_ROOT__/scripts/run_weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>ProcessType</key>
  <string>Background</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PROJECT_ROOT__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/logs/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/logs/launchd.stderr.log</string>

  <key>ThrottleInterval</key>
  <integer>60</integer>
</dict>
</plist>
